# supervisor-thumbnailer.conf
[program:thumbnailer]
command=/home/eventials/thumbnailer/current/thumbnailer.sh
directory=/home/eventials/thumbnailer/current/
user=eventials
autostart=true
autorestart=true
